# foundry/provisioning/roles/raspbian/files/systemd/update.path
[Unit]
AssertPathExists=/home/pi/updates
Wants=update.service

[Path]
PathChanged=/home/pi/updates

[Install]
WantedBy=multi-user.target
