# foundry/provisioning/roles/raspbian/files/systemd/upload.path
[Unit]
AssertPathExists=/home/pi/storage
Wants=update.service

[Path]
PathChanged=/home/pi/storage/studentcode.py

[Install]
WantedBy=multi-user.target
