# DevOps/fabric/systemd-units/fabric-update.path
[Path]
Unit=fabric-update.service
MakeDirectory=true
DirectoryNotEmpty=/home/fabric/updates

[Install]
WantedBy=multi-user.target
