# frankfurter/resources/memcached.conf
start on starting runtime
setuid memcache

script
  exec memcached -p 12357
end script
