# DevOps/frankfurter/resources/runtime.service
[Unit]
Description=runtime service
Wants=update.service
After=multi-user.target update.service

[Service]
Type=simple
ExecStart=/home/ubuntu/bin/runtime.sh

[Install]
WantedBy=multi-user.target
